# doppler/config/doppler-agent.upstart
#!upstart
description "Doppler metrics collection agent"
author      "Doppler <agent@doppler.io>"

start on (filesystem and net-device-up IFACE=lo)
stop on shutdown

respawn
respawn limit 99 5

script
    exec /usr/bin/doppler-agent.py 2>> /tmp/doppler.out 1>> /tmp/doppler.err
end script
